# packaging/systemd/eightyeightyeight.service
# NOTE: This systemd service file is a stub for demonstration purposes.
# In a production environment, you would need to adjust paths, users, 
# and ensure proper environment variable management.

[Unit]
Description=Eightyeightyeight Video Capture Service
After=network.target

[Service]
ExecStart=/usr/local/bin/eightyeightyeight record
Restart=on-failure
User=root
Group=root
EnvironmentFile=/etc/eightyeightyeight/config.toml

[Install]
WantedBy=multi-user.target
