# packaging/systemd/eightyeightyeight-stream.service
# NOTE: This systemd service file is a stub for demonstration purposes.
# In a production environment, you would need to adjust paths, users, 
# and ensure proper environment variable management.

[Unit]
Description=Eightyeightyeight Video Streaming Service
After=network.target

[Service]
# Assumes config is placed at /etc/eightyeightyeight/config.toml
# You might want to override DEST and PORT via drop-in files or environment variables
ExecStart=/usr/local/bin/eightyeightyeight stream --config /etc/eightyeightyeight/config.toml --dest 127.0.0.1 --port 8088
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target
