# packaging/systemd/eightyeightyeight-record.service
# NOTE: This systemd service file is a stub for demonstration purposes.
# In a production environment, you would need to adjust paths, users, 
# and ensure proper environment variable management.

[Unit]
Description=Eightyeightyeight Video Recording Service
After=network.target

[Service]
# Assumes config is placed at /etc/eightyeightyeight/config.toml
ExecStart=/usr/local/bin/eightyeightyeight record --config /etc/eightyeightyeight/config.toml
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target